# upstart/aerokubot.conf
description "Aerokubot telegram bot"

# start in normal runlevels when disks are mounted and networking is available
start on runlevel [2345]

# stop on shutdown/halt, single-user mode and reboot
stop on runlevel [016]

#pre-start exec set -o allexport; source /opt/aerokubot/.env; set +o allexport
exec sudo -u aerokubot /opt/aerokubot/aerokubot --token ${TELEGRAM_TOKEN}
